# systemd/github-agent-bridge-monitor.service
[Unit]
Description=GitHub Agent Bridge monitor
After=github-agent-bridge.service

[Service]
Type=oneshot
EnvironmentFile=-%h/.config/github-agent-bridge/env
ExecStart=%h/.local/bin/github-agent-bridge-monitor-alert
