# systemd/github-agent-bridge-reader.service
[Unit]
Description=GitHub Agent Bridge IMAP reader
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=GITHUB_AGENT_BRIDGE_IMAP_HOST=imap.gmail.com
Environment=GITHUB_AGENT_BRIDGE_IMAP_PORT=993
Environment=GITHUB_AGENT_BRIDGE_MAILBOX=INBOX
Environment=GITHUB_AGENT_BRIDGE_MARK_SEEN=
EnvironmentFile=-%h/.config/github-agent-bridge/env
ExecStart=%h/.local/bin/github-agent-bridge-reader-run
