# systemd/github-agent-bridge.service
[Unit]
Description=GitHub Agent Bridge
After=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/github-agent-bridge
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
